# deploy/proposal-reviewer.service
[Unit]
Description=AI Proposal Reviewer - Groq API Service
After=network.target

[Service]
Type=notify
User=viona
Group=viona
WorkingDirectory=/opt/proposal-reviewer
Environment="PATH=/opt/proposal-reviewer/venv/bin"
ExecStart=/opt/proposal-reviewer/venv/bin/gunicorn app.utama:aplikasi \
    --workers 4 \
    --worker-class uvicorn.workers.UvicornWorker \
    --bind 127.0.0.1:8000 \
    --access-logfile /var/log/proposal-reviewer/access.log \
    --error-logfile /var/log/proposal-reviewer/error.log \
    --log-level info
ExecReload=/bin/kill -s HUP $MAINPID
KillMode=mixed
TimeoutStopSec=5
PrivateTmp=true

# Restart policy
Restart=on-failure
RestartSec=10s

[Install]
WantedBy=multi-user.target
